# debian/quantum-server.prerm
#!/bin/sh

set -e

case $1 in
        remove|purge)
                invoke-rc.d quantum-server stop
esac

#DEBHELPER#
